# .config/systemd/user/backup.timer
[Timer]
OnUnitActiveSec=1d
OnStartupSec=0s

[Install]
WantedBy=default.target
