# .config/systemd/user/backup.service
[Unit]
Description=Run a local backup
OnFailure=reschedule-backup.timer

[Service]
ExecStart=%E/systemd/user/backup.fish
