# .config/systemd/user/reschedule-backup.timer
[Timer]
Unit=backup.service
OnActiveSec=1h
RemainAfterElapse=false
