# .config/systemd/user/reschedule-backup.service
[Unit]
Description=Reschedule a failed backup.

[Service]
ExecStart=%E/systemd/user/reschedule-backup.fish
