# .config/systemd/user/scheduled-backup@.timer
[Timer]
Unit=backup.service
OnCalendar=@%i
Persistent=true
RemainAfterElapse=false
